# root/etc/s6-overlay/s6-rc.d/init-data/run
#!/usr/bin/with-contenv bash

mkdir -p /data/acme.sh /data/cert /data/derper
lsiown  -R abc:abc /data
